# cron.d/robin
* * * * * root /sbin/start reddit-job-robin_prompt_for_voting
* * * * * root /sbin/start reddit-job-robin_reap_ripe_rooms

